# build/package/linux_agent/ipk/pbs-plus-updater
#!/bin/sh /etc/rc.common

START=99
STOP=10

USE_PROCD=1

PROG=/usr/bin/pbs-plus-updater
PIDFILE=/run/proxmox-backup/pbs-plus-updater.pid

start_service() {
    # Create run directory if it doesn't exist
    mkdir -p /run/proxmox-backup
    
    procd_open_instance
    procd_set_param command "$PROG"
    procd_set_param pidfile "$PIDFILE"
    procd_set_param user root
    procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-5}
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

reload_service() {
    if [ -f "$PIDFILE" ]; then
        kill -HUP $(cat "$PIDFILE")
    fi
}

service_triggers() {
    procd_add_reload_trigger "pbs-plus-updater"
}
