# build/package/linux_agent/deb/pbs-plus-updater.service
[Unit]
Description=PBS Plus Agent Updater
Wants=network-online.target
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/pbs-plus-updater
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/pbs-plus-agent/pbs-plus-updater.pid
Restart=on-failure
RestartSec=30

# Working directory
WorkingDirectory=/var/lib/pbs-plus-agent

# Security settings - updater needs more privileges
NoNewPrivileges=false # Allows the process to gain new privileges (necessary for root)
PrivateTmp=false      # dpkg might need access to the global /tmp
ProtectSystem=false   # Allows dpkg to write to system directories like /usr, /etc
ProtectHome=true      # Generally safe; dpkg usually doesn't modify /home. Set to false if necessary.

# Environment
Environment=HOME=/root # Standard home directory for root processes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pbs-plus-updater

[Install]
WantedBy=multi-user.target
