# build/package/linux_agent/deb/postinst
#!/bin/bash
set -e

# Create dedicated user and group if they don't exist
if ! getent group pbs-plus >/dev/null 2>&1; then
    groupadd -r pbs-plus
fi

if ! getent passwd pbs-plus >/dev/null 2>&1; then
    useradd -r -g pbs-plus -s /bin/false \
        -d /var/lib/pbs-plus-agent \
        -c "PBS Plus Agent" pbs-plus
fi

# Set executable permissions
chmod +x /usr/bin/pbs-plus-agent
chmod +x /usr/bin/pbs-plus-updater

# Create necessary directories with proper ownership
mkdir -p /var/lib/pbs-plus-agent
mkdir -p /var/log/pbs-plus-agent
mkdir -p /run/proxmox-backup
mkdir -p /etc/pbs-plus-agent

# Set ownership and permissions
chown pbs-plus:pbs-plus /var/lib/pbs-plus-agent
chown pbs-plus:pbs-plus /var/log/pbs-plus-agent
chown pbs-plus:pbs-plus /run/proxmox-backup
chown pbs-plus:pbs-plus /etc/pbs-plus-agent

chmod 750 /var/lib/pbs-plus-agent
chmod 750 /var/log/pbs-plus-agent
chmod 755 /run/proxmox-backup
chmod 750 /etc/pbs-plus-agent

# Set binary ownership
chown root:pbs-plus /usr/bin/pbs-plus-agent
chown root:pbs-plus /usr/bin/pbs-plus-updater
chmod 755 /usr/bin/pbs-plus-agent
chmod 755 /usr/bin/pbs-plus-updater

# Add pbs-plus to backup group if it exists
if getent group backup >/dev/null 2>&1; then
    usermod -a -G backup pbs-plus
fi

# Set file capabilities instead of using sudo
setcap 'cap_dac_read_search,cap_fowner,cap_dac_override+ep' /usr/bin/pbs-plus-agent
setcap 'cap_dac_read_search,cap_fowner,cap_dac_override,cap_setuid,cap_setgid,cap_sys_admin+ep' /usr/bin/pbs-plus-updater

systemctl daemon-reload

# Enable the services to start on boot
systemctl enable pbs-plus-agent.service
systemctl enable pbs-plus-updater.service

# Start the services
systemctl restart pbs-plus-agent.service
systemctl restart pbs-plus-updater.service

exit 0
