# build/package/linux_agent/apk/pbs-plus-agent
#!/sbin/openrc-run

name="PBS Plus Agent"
description="PBS Plus Agent for backup operations"

command="/usr/bin/pbs-plus-agent"
command_user="pbsplus:pbsplus"
pidfile="/run/pbs-plus-agent/pbs-plus-agent.pid"
command_background="yes"

# Working directory
directory="/var/lib/pbs-plus-agent"

# Output redirection
output_log="/var/log/pbs-plus-agent/agent.log"
error_log="/var/log/pbs-plus-agent/agent.error.log"

depend() {
    need net
    after firewall
}

start_pre() {
    # Create run directory if it doesn't exist
    checkpath --directory --owner pbsplus:pbsplus --mode 0755 /run/pbs-plus-agent
    
    # Ensure log directory exists
    checkpath --directory --owner pbsplus:pbsplus --mode 0750 /var/log/pbs-plus-agent
    
    # Ensure working directory exists
    checkpath --directory --owner pbsplus:pbsplus --mode 0750 /var/lib/pbs-plus-agent

    # Ensure config directory exists
    checkpath --directory --owner pbsplus:pbsplus --mode 0750 /etc/pbs-plus-agent
}

reload() {
    ebegin "Reloading $name"
    if [ -f "$pidfile" ]; then
        kill -HUP $(cat "$pidfile")
        eend $?
    else
        eend 1 "PID file not found"
    fi
}

stop_post() {
    # Clean up PID file if it exists
    [ -f "$pidfile" ] && rm -f "$pidfile"
}
