# build/package/linux_agent/apk/postinst
#!/bin/sh
set -e

chmod +x /usr/bin/pbs-plus-agent

# Create run directory if it doesn't exist
mkdir -p /run/proxmox-backup

# Add the service to the default runlevel
rc-update add pbs-plus-agent default

# Start the service
rc-service pbs-plus-agent start

exit 0
