# build/package/linux_agent/rpm/pbs-plus-updater.service
[Unit]
Description=PBS Plus Agent Updater
Wants=network-online.target
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/pbs-plus-updater
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/pbs-plus-agent/pbs-plus-updater.pid
Restart=on-failure
RestartSec=30
User=pbsplus
Group=pbsplus

# Working directory
WorkingDirectory=/var/lib/pbs-plus-agent

# Security settings - updater needs more privileges
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=true
ReadWritePaths=/etc/pbs-plus-agent /var/lib/pbs-plus-agent /var/log/pbs-plus-agent /run/pbs-plus-agent /usr/bin

# Environment
Environment=HOME=/var/lib/pbs-plus-agent
Environment=USER=pbsplus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pbs-plus-updater

[Install]
WantedBy=multi-user.target
