# build/package/linux_agent/rpm/pbs-plus-updater.service
[Unit]
Description=PBS Plus Agent Updater
Wants=network-online.target
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/pbs-plus-updater
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/pbs-plus-agent/pbs-plus-updater.pid
Restart=on-failure
RestartSec=30

# Working directory
WorkingDirectory=/var/lib/pbs-plus-agent

# Security settings - updater needs more privileges
NoNewPrivileges=false # Allows the process to gain new privileges (necessary for root)
PrivateTmp=false      # Package managers might need access to the global /tmp
ProtectSystem=false   # Allows rpm/yum/dnf to write to system directories
ProtectHome=true      # Generally safe; package managers usually don't modify /home.

# Environment
Environment=HOME=/root # Standard home directory for root processes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pbs-plus-updater

[Install]
WantedBy=multi-user.target
